# examples/service/units/user/swish-service-example.service.in
[Unit]
Description=A simple user-mode swish service

[Service]
Type=simple
ExecStart=%SED_PWD%/swish-service-example /SERVICE
WorkingDirectory=%SED_PWD%
Restart=always

[Install]
WantedBy=default.target
